# systemd/oled-menu.service
[Unit]
Description=STREAMER OLED menu service (user)
After=default.target mpd.service

[Service]
Type=simple
Environment="ENC_A=24" "ENC_B=23" "BTN=13" "DIM_TIMEOUT=30" "INVERT_ROT=0"
WorkingDirectory=%h/streamer
ExecStart=%h/streamer/venv/bin/python %h/streamer/oled/menu_service.py
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
